# systemd/openclaw-ai-quant-prune-runtime-logs.service.example
[Unit]
Description=OpenClaw AI Quant Prune Runtime Logs (SQLite)

[Service]
Type=oneshot
WorkingDirectory=$PROJECT_DIR
EnvironmentFile=-%h/.config/openclaw/ai-quant-live.env
ExecStart=$PROJECT_DIR/.venv/bin/python3 -u tools/prune_runtime_logs.py

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes

[Install]
WantedBy=default.target

