# systemd/openclaw-ai-quant-prune-runtime-logs.service.example
[Unit]
Description=OpenClaw AI Quant Prune Runtime Logs (SQLite)

[Service]
Type=oneshot
WorkingDirectory=$PROJECT_DIR
EnvironmentFile=-%h/.config/openclaw/ai-quant-live.env
ExecStart=$PROJECT_DIR/.venv/bin/python3 -u tools/prune_runtime_logs.py

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadWritePaths=$PROJECT_DIR

[Install]
WantedBy=default.target

