# systemd/openclaw-ai-quant-monitor.service.example
[Unit]
Description=OpenClaw AI Quant Monitor Dashboard
After=network-online.target openclaw-ai-quant-ws-sidecar.service
Wants=network-online.target openclaw-ai-quant-ws-sidecar.service
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
WorkingDirectory=$PROJECT_DIR
EnvironmentFile=-%h/.config/openclaw/ai-quant-monitor.env
ExecStart=$PROJECT_DIR/.venv/bin/python3 -u monitor/server.py
Restart=always
RestartSec=5
KillMode=process

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallFilter=@system-service
ReadWritePaths=$PROJECT_DIR

# Resource limits
MemoryMax=256M
TasksMax=256
TimeoutStopSec=15

[Install]
WantedBy=default.target
