# systemd/openclaw-ai-quant-prune-runtime-logs.timer.example
[Unit]
Description=OpenClaw AI Quant Prune Runtime Logs Timer

[Timer]
OnCalendar=*-*-* 01:15:00 UTC
Persistent=true
Unit=openclaw-ai-quant-prune-runtime-logs.service

[Install]
WantedBy=timers.target

