# systemd/openclaw-ai-quant-hub.service.example
[Unit]
Description=OpenClaw AI Quant Hub
After=network-online.target openclaw-ai-quant-ws-sidecar.service
Wants=network-online.target openclaw-ai-quant-ws-sidecar.service
Requires=openclaw-ai-quant-ws-sidecar.service

[Service]
Type=simple
WorkingDirectory=$PROJECT_DIR
EnvironmentFile=-%h/.config/openclaw/ai-quant-monitor-v8.env
Environment=AI_QUANT_WS_SIDECAR_SOCK=%t/openclaw-ai-quant-ws.sock
Environment=AI_QUANT_CANDLES_DB_DIR=$PROJECT_DIR/candles_dbs
Environment=AIQ_ROOT=$PROJECT_DIR
Environment=RUST_LOG=info
Environment=PATH=%h/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# Enable Hub manual-trade writes in the Hub process itself.
# The Hub does not require AI_QUANT_LIVE_ENABLE / AI_QUANT_LIVE_CONFIRM for
# manual-trade routes; use AI_QUANT_HARD_KILL_SWITCH=1 to fail closed.
# Read-only Hub routes may bypass AIQ_MONITOR_TOKEN for trusted LAN/Tailscale
# clients. Loopback alone does not bypass read auth; use a Tailscale
# HTTP/HTTPS Serve-style front layer (not raw TCP localhost forwarding) when
# you need Tailscale identity-header bypass.
# Configure these tokens in the environment file:
# - AIQ_MONITOR_TOKEN for viewer/read access
# - AIQ_MONITOR_EDITOR_TOKEN for config saves and live request creation
# - AIQ_MONITOR_APPROVER_TOKEN for live approval / rejection
# - AIQ_MONITOR_ADMIN_TOKEN for privileged diagnostics and system actions
# If an operator explicitly trusts bare localhost forwarding on this machine,
# they can opt in with:
# Environment=AIQ_MONITOR_TRUST_LOOPBACK_READ=1
# Environment=AIQ_MONITOR_TRUST_LOOPBACK_ADMIN=1
# Keep AIQ_MONITOR_TOKEN / AIQ_MONITOR_ADMIN_TOKEN configured in the Hub
# environment. These flags only bypass presenting credentials for loopback
# requests.
Environment=AIQ_MANUAL_TRADE_ENABLE=true
ExecStart=$PROJECT_DIR/hub/target/release/aiq-hub
Restart=always
RestartSec=5
KillMode=process

[Install]
WantedBy=default.target
