# systemd/openclaw-ai-quant-prune-runtime-logs-v8.service.example
[Unit]
Description=OpenClaw AI Quant Prune Runtime Logs (SQLite) (v8)

[Service]
Type=oneshot
WorkingDirectory=$PROJECT_DIR
EnvironmentFile=-%h/.config/openclaw/ai-quant-v8.env
ExecStart=$PROJECT_DIR/.venv/bin/python3 -u tools/prune_runtime_logs.py

[Install]
WantedBy=default.target

