# systemd/openclaw-ai-quant-hub.service.example
[Unit]
Description=OpenClaw AI Quant Hub
After=network-online.target openclaw-ai-quant-ws-sidecar.service
Wants=network-online.target openclaw-ai-quant-ws-sidecar.service
Requires=openclaw-ai-quant-ws-sidecar.service

[Service]
Type=simple
WorkingDirectory=$PROJECT_DIR
EnvironmentFile=-%h/.config/openclaw/ai-quant-monitor-v8.env
Environment=AI_QUANT_WS_SIDECAR_SOCK=%t/openclaw-ai-quant-ws.sock
Environment=AI_QUANT_CANDLES_DB_DIR=$PROJECT_DIR/candles_dbs
Environment=AIQ_ROOT=$PROJECT_DIR
Environment=RUST_LOG=info
Environment=PATH=%h/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# Enable Hub manual-trade writes in the Hub process itself.
# The Hub does not require AI_QUANT_LIVE_ENABLE / AI_QUANT_LIVE_CONFIRM for
# manual-trade routes; use AI_QUANT_HARD_KILL_SWITCH=1 to fail closed.
# Read-only Hub routes may bypass AIQ_MONITOR_TOKEN for trusted LAN/Tailscale
# clients. Loopback alone does not bypass read auth; admin/manual-trade routes
# still require AIQ_MONITOR_ADMIN_TOKEN.
Environment=AIQ_MANUAL_TRADE_ENABLE=true
ExecStart=$PROJECT_DIR/hub/target/release/aiq-hub
Restart=always
RestartSec=5
KillMode=process

[Install]
WantedBy=default.target
