# systemd/openclaw-ai-quant-factory-v8.timer.example
[Unit]
Description=OpenClaw AI Quant Nightly Strategy Factory Timer (v8)

[Timer]
OnCalendar=*-*-* 00:50:00 UTC
Persistent=true
Unit=openclaw-ai-quant-factory-v8.service

[Install]
WantedBy=timers.target

