# systemd/openclaw-ai-quant-trader-v8-livepaper.service.example
[Unit]
Description=OpenClaw AI Quant Rust Paper Lane (v8 livepaper)
After=network-online.target openclaw-ai-quant-ws-sidecar.service
Wants=network-online.target
Requires=openclaw-ai-quant-ws-sidecar.service

[Service]
Type=simple
WorkingDirectory=$PROJECT_DIR
EnvironmentFile=-%h/.config/openclaw/ai-quant-universe-v8.env
EnvironmentFile=-%h/.config/openclaw/ai-quant-v8.env
EnvironmentFile=-%h/.config/openclaw/ai-quant-trader-v8-livepaper.env
ExecStart=$PROJECT_DIR/scripts/run_paper_lane.sh livepaper
Restart=always
RestartSec=5
KillMode=process

Environment=AI_QUANT_RUNTIME_BIN=$PROJECT_DIR/target/release/aiq-runtime

[Install]
WantedBy=default.target
