# systemd/openclaw-ai-quant-factory-v8.service.example
[Unit]
Description=OpenClaw AI Quant Nightly Strategy Factory Cycle (v8 parallel instance)
After=network-online.target openclaw-ai-quant-ws-sidecar-v8.service
Wants=network-online.target
Requires=openclaw-ai-quant-ws-sidecar-v8.service
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=oneshot
TimeoutStartSec=7200
WorkingDirectory=$PROJECT_DIR
EnvironmentFile=-%h/.config/openclaw/ai-quant-universe-v8.env
EnvironmentFile=-%h/.config/openclaw/ai-quant-v8.env

# Use a v8-scoped backtester binary (do not overwrite production).
Environment=MEI_BACKTESTER_BIN=$PROJECT_DIR/backtester/dist/mei-backtester-gpu

# WSL2 CUDA: ensure libcuda can be found by cudarc.
# On native Linux, you can usually remove this line.
Environment=LD_LIBRARY_PATH=/usr/lib/wsl/lib
Environment=AI_QUANT_FACTORY_HEALTH_SERVICES=openclaw-ai-quant-trader-v8-paper1,openclaw-ai-quant-trader-v8-paper2,openclaw-ai-quant-trader-v8-paper3,openclaw-ai-quant-trader-v8-livepaper

# Prevent overlapping runs.
ExecStart=/usr/bin/flock -n %t/openclaw-ai-quant-factory-v8.lock /usr/bin/env bash -lc "cd \"$PROJECT_DIR\" && \"$PROJECT_DIR/.venv/bin/python3\" -u tools/factory_cycle.py --profile daily --strategy-mode primary --gpu --tpe --walk-forward --slippage-stress --concentration-checks --sensitivity-checks --max-age-fail-hours 24 --funding-max-stale-symbols 1 --service openclaw-ai-quant-trader-v8-paper1 --yaml-path \"$PROJECT_DIR/config/strategy_overrides.paper1.yaml\" --candidate-count 3 --candidate-services openclaw-ai-quant-trader-v8-paper1,openclaw-ai-quant-trader-v8-paper2,openclaw-ai-quant-trader-v8-paper3 --candidate-yaml-paths \"$PROJECT_DIR/config/strategy_overrides.paper1.yaml,$PROJECT_DIR/config/strategy_overrides.paper2.yaml,$PROJECT_DIR/config/strategy_overrides.paper3.yaml\" --enable-livepaper-promotion --livepaper-service openclaw-ai-quant-trader-v8-livepaper --livepaper-yaml-path \"$PROJECT_DIR/config/strategy_overrides.livepaper.yaml\" --ws-service openclaw-ai-quant-ws-sidecar-v8 --discord-target \"${AI_QUANT_FACTORY_DISCORD_TARGET}\""

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallFilter=@system-service
ReadWritePaths=$PROJECT_DIR

# Resource limits
MemoryMax=1G
TasksMax=256
TimeoutStopSec=15

[Install]
WantedBy=default.target
