# systemd/openclaw-ai-quant-prune-runtime-logs.service.example
[Unit]
Description=OpenClaw AI Quant Prune Runtime Logs (SQLite)

[Service]
Type=oneshot
WorkingDirectory=$PROJECT_DIR
EnvironmentFile=-%h/.config/openclaw/ai-quant-live.env
ExecStart=$PROJECT_DIR/.venv/bin/python3 -u tools/prune_runtime_logs.py

# Security hardening
ProtectSystem=strict
NoNewPrivileges=yes
PrivateTmp=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadWritePaths=/home/fol2hk/openclaw-plugins/ai_quant

[Install]
WantedBy=default.target

