# systemd/openclaw-ai-quant-monitor-v8.service.example
[Unit]
Description=OpenClaw AI Quant Monitor Dashboard (v8 parallel instance)
After=network-online.target openclaw-ai-quant-ws-sidecar-v8.service
Wants=network-online.target openclaw-ai-quant-ws-sidecar-v8.service

[Service]
Type=simple
WorkingDirectory=$PROJECT_DIR
EnvironmentFile=-%h/.config/openclaw/ai-quant-monitor-v8.env
Environment=AI_QUANT_WS_SIDECAR_SOCK=%t/openclaw-ai-quant-ws-v8.sock
ExecStart=$PROJECT_DIR/.venv/bin/python3 -u monitor/server.py
Restart=always
RestartSec=5
KillMode=process

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadWritePaths=$PROJECT_DIR

# Resource limits
MemoryMax=256M
TasksMax=256
TimeoutStopSec=15
StartLimitBurst=5
StartLimitIntervalSec=60

[Install]
WantedBy=default.target
