# systemd/openclaw-ai-quant-trader-v8-paper2.service.example
[Unit]
Description=OpenClaw AI Quant Paper Trader Daemon (v8 paper2)
After=network-online.target openclaw-ai-quant-ws-sidecar.service
Wants=network-online.target
Requires=openclaw-ai-quant-ws-sidecar.service

[Service]
Type=simple
WorkingDirectory=$PROJECT_DIR
EnvironmentFile=-%h/.config/openclaw/ai-quant-universe-v8.env
EnvironmentFile=-%h/.config/openclaw/ai-quant-v8.env
EnvironmentFile=-%h/.config/openclaw/ai-quant-trader-v8-paper2.env
ExecStartPre=/usr/bin/test -x $PROJECT_DIR/target/release/aiq-runtime
ExecStartPre=/usr/bin/mkdir -p $PROJECT_DIR/artifacts/state
ExecStartPre=/usr/bin/touch $PROJECT_DIR/artifacts/state/paper_watchlist_paper2.txt
ExecStart=$PROJECT_DIR/target/release/aiq-runtime paper lane daemon --lane paper2 --project-dir $PROJECT_DIR --symbols-file $PROJECT_DIR/artifacts/state/paper_watchlist_paper2.txt --watch-symbols-file --lookback-bars 200
Restart=always
RestartSec=5
KillMode=process

# Build the Rust runtime before enabling this lane:
#   cargo build --release -p aiq-runtime

Environment=RUST_BACKTRACE=1

[Install]
WantedBy=default.target
