# systemd/openclaw-ai-quant-trader.service.example
[Unit]
Description=OpenClaw AI Quant Paper Trader Daemon
After=network-online.target openclaw-ai-quant-ws-sidecar.service
Wants=network-online.target
Requires=openclaw-ai-quant-ws-sidecar.service

[Service]
Type=simple
WorkingDirectory=$PROJECT_DIR
EnvironmentFile=-%h/.config/openclaw/ai-quant-universe.env
ExecStart=$PROJECT_DIR/.venv/bin/python3 -u -m engine.daemon
Restart=always
RestartSec=5
KillMode=process
Environment=AI_QUANT_MODE=paper
Environment=AI_QUANT_LOCK_PATH=$PROJECT_DIR/ai_quant_paper.lock
Environment=AI_QUANT_INTERVAL=1h
Environment=AI_QUANT_LOOKBACK_BARS=200
Environment=AI_QUANT_NEUTRAL_AUDIT_SAMPLE_EVERY_S=600
Environment=AI_QUANT_NEUTRAL_AUDIT_SAMPLE_SYMBOLS=5
Environment=AI_QUANT_DEBUG_GATES_SYMBOLS=BTC
Environment=AI_QUANT_DEBUG_GATES_EVERY_S=30
Environment=AI_QUANT_WS_ENABLE_BBO=0
Environment=AI_QUANT_WS_STALE_BBO_S=0
Environment=AI_QUANT_WS_RESTART_COOLDOWN_S=120
Environment=AI_QUANT_WS_RESTART_MAX_IN_WINDOW=8
Environment=AI_QUANT_WS_SOURCE=sidecar
Environment=AI_QUANT_WS_SIDECAR_SOCK=%t/openclaw-ai-quant-ws.sock
Environment=AI_QUANT_MARKET_DB_PATH=$PROJECT_DIR/market_data.db
Environment=AI_QUANT_CANDLES_SOURCE=sidecar
Environment=AI_QUANT_CANDLES_DB_DIR=$PROJECT_DIR/candles_dbs
Environment=AI_QUANT_REST_CANDLE_BACKFILL=0
Environment=AI_QUANT_REST_ENABLE=0
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
